# bootc/arch-base/rootfs/usr/lib/systemd/system/systemd-journal-flush.service.d/shutdown-fix.conf
# Workaround for https://github.com/systemd/systemd/issues/40199
[Unit]
Before=umount.target
Conflicts=umount.target
