# bootc/arch-common-config/rootfs/usr/lib/systemd/system/nix.mount
[Unit]
Description=Bind mount for persistent nix store
Requires=var.mount
After=var.mount

[Mount]
What=/var/nix
Where=/nix
Options=bind
Type=none