# bootc/arch-common-config/rootfs/usr/lib/systemd/system/provision.target
[Unit]
Description=Provision procedures for the system image